# files/etc/uci-defaults/90_squid_profiles
#!/bin/sh
# Create default squid_profiles UCI configuration if missing

CONFIG_FILE=/etc/config/squid_profiles

backup_and_initialize() {
    # Backup the existing Squid configuration if present and create
    # a minimal default skeleton.  This function renames the file
    # /etc/squid/squid.conf to include today's date (YYYYMMDD) and
    # writes a new configuration skeleton that is compatible with
    # luci-app-squid-profiles.  It also prepares the domain and map
    # directories used by the plugin.

    SQUID_DIR="/etc/squid"
    SQUID_CONF="${SQUID_DIR}/squid.conf"
    today="$(date +%Y%m%d)"
    if [ -f "$SQUID_CONF" ]; then
        # Only back up if we haven't already done so on this date.
        backup_path="${SQUID_CONF}.${today}.bak"
        if [ ! -f "$backup_path" ]; then
            mv "$SQUID_CONF" "$backup_path"
            echo "Backed up existing Squid configuration to $backup_path" >&2
        fi
    fi

    # Create the directory structure for domains and maps if missing
    mkdir -p "$SQUID_DIR/domains" "$SQUID_DIR/maps"

    # Write a minimal default squid.conf if it does not exist.  This
    # skeleton contains basic settings and placeholders; it will not
    # override an existing file if one remains after backup.
    if [ ! -f "$SQUID_CONF" ]; then
        cat > "$SQUID_CONF" <<'EOCONF'
#
# Default Squid configuration generated by luci-app-squid-profiles
# This file provides a minimal skeleton suitable for use with the
# dynamic ACLs managed by the LuCI interface.  It may be customised
# further as needed.  Existing configurations are preserved via a
# time-stamped backup.

visible_hostname openwrt-squid-profiles

http_port 3128

# Safe ports
acl SSL_ports port 443
acl Safe_ports port 80
acl Safe_ports port 443
acl CONNECT method CONNECT

# Deny requests to unknown ports
http_access deny !Safe_ports
http_access deny CONNECT !SSL_ports

# Placeholder for dynamically generated profile rules
# Profile ACLs and http_access rules are inserted by luci-app-squid-profiles

# Deny everything else by default
http_access deny all

# Disable caching on systems with limited storage
cache deny all
cache_mem 8 MB
maximum_object_size 0 KB
pinger_enable off

# Logging to temporary files
access_log stdio:/tmp/squid_access.log squid
cache_log /tmp/squid_cache.log
logfile_rotate 3

# Retain default refresh patterns
refresh_pattern ^ftp:       1440    20%     10080
refresh_pattern ^gopher:    1440    0%      1440
refresh_pattern -i (/cgi-bin/|\?) 0 0%      0
refresh_pattern .           0       20%     4320

# End of default configuration
EOCONF
        echo "Created default Squid configuration skeleton at $SQUID_CONF" >&2
    fi
}

create_default_network() {
    # Create at least one network section covering the typical LAN if
    # the squid_profiles configuration does not exist.  Each network
    # section holds a subnet in CIDR notation and the corresponding
    # VLAN identifier.  Additional networks can be added later through
    # the LuCI UI.
    if ! uci show squid_profiles >/dev/null 2>&1; then
        uci batch <<'EOF'
add squid_profiles network
set squid_profiles.@network[-1].subnet='192.168.1.0/24'
set squid_profiles.@network[-1].vlan='1'
commit squid_profiles
EOF
    fi
}

main() {
    # Perform backup and initialization of Squid and squid_profiles
    # configuration files.  Afterwards, create a default network if
    # the squid_profiles UCI config is missing.
    backup_and_initialize
    create_default_network
}

main "$@"

exit 0