# openwrt-squid-profile-ui/files/etc/init.d/squid-profiles
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

HELPER="/usr/libexec/squid-profiles"
CONFIG="squid_profiles"

is_enabled() {
	[ "$(uci -q get "$CONFIG.core.enabled" 2>/dev/null || echo 1)" = "1" ]
}

start_service() {
	if is_enabled; then
		[ -x "$HELPER" ] && "$HELPER" init
	else
		echo "squid-profiles disabled; Squid configuration left untouched"
	fi
}

reload_service() {
	if is_enabled; then
		[ -x "$HELPER" ] && "$HELPER" apply
	else
		echo "squid-profiles disabled; refusing to apply"
	fi
}

service_triggers() {
	procd_add_reload_trigger "$CONFIG"
}

boot() {
	start_service
}
