# openwrt-squid-profile-ui/files/etc/uci-defaults/90_squid_profiles
#!/bin/sh

set -eu

HELPER="/usr/libexec/squid-profiles"

if [ ! -x "$HELPER" ]; then
	echo "squid-profiles helper is missing: $HELPER" >&2
	exit 1
fi

"$HELPER" init

if ! uci -q get squid_profiles.core >/dev/null 2>&1; then
	uci batch <<'EOUCI'
set squid_profiles.core=globals
set squid_profiles.core.enabled='1'
commit squid_profiles
EOUCI
fi

exit 0
